# milestone1/raspberrypi/RoboRetriever.service
[Unit]
Description=My Script Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/myscript.sh
Restart=always
RestartSec=5
User=your-username
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target